# rubigo-react/rubigo-react.service
[Unit]
Description=Rubigo React (Next.js on Bun)
After=network.target

[Service]
Type=simple
User=kstewart
WorkingDirectory=/home/kstewart/development/rubigo/production/rubigo-react/current
ExecStart=/home/kstewart/.bun/bin/bun run start -- -p 4430
Environment=NODE_ENV=production
Environment=DATABASE_URL=file:../data/rubigo.db
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
