# rubigo-react.service
[Unit]
Description=Rubigo React (Next.js on Bun)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/kstewart/development/rubigo/runner/_work/rubigo/rubigo/production/rubigo-react/current
ExecStart=/usr/local/bin/bun run start -- -p 4430
Environment=NODE_ENV=production
Environment=DATABASE_URL=/home/kstewart/development/rubigo/runner/_work/rubigo/rubigo/production/rubigo-react/data/rubigo.db
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target